# backend/docker/cron/cron.d/yii-cron
* * * * * /usr/local/bin/php /app/yii queue/run >> /var/log/cron.log 2>&1
